# deploy/systemd/dextrabot-db-retention.service
[Unit]
Description=Dextrabot high-volume DB retention cleanup
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=dextrabot
Group=dextrabot
WorkingDirectory=__DEXTRABOT_ROOT__
EnvironmentFile=/etc/dextrabot/dextrabot.env
Environment=RETENTION_DAYS=3
Environment=BATCH_SIZE=50000
ExecStart=__DEXTRABOT_ROOT__/scripts/prune_high_volume_db_tables.sh --apply --vacuum

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=false
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
