# deploy/systemd/dextrabot-db-retention.timer
[Unit]
Description=Run Dextrabot high-volume DB retention cleanup hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Persistent=true
Unit=dextrabot-db-retention.service

[Install]
WantedBy=timers.target
